# usbdumper/mount.rules
ENV{DEVTYPE}=="usb_device", ACTION=="add", SUBSYSTEM=="usb", RUN+="/bin/dumpusb /dev/sdb &"
